# supervisord.conf
[program:hookbox]
command={{PROJECT_DIR}}/hookbox.sh
